# Ansible/services/clamav-scan@home.timer
[Unit]
Description=Daily ClamAV scan of /home

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
